# recipes-multimedia/alsa/audio-detect/audio-detect.service
[Unit]
Description=Audio Hardware Detection and Validation
After=sound.target
Before=mpd-auris.service upmpdcli.service shairport-sync.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/audio-detect.sh
RemainAfterExit=yes

[Install]
WantedBy=sound.target
