# recipes-multimedia/shairport-sync/files/shairport-sync.service
[Unit]
Description=Shairport Sync AirPlay Audio Receiver
After=sound.target network.target avahi-daemon.service
Requires=avahi-daemon.service

[Service]
Type=simple
User=shairport-sync
Group=audio

# Restrict to CPU cores 0-2 (not using isolated CPU 3)
CPUAffinity=0-2

ExecStart=/usr/bin/shairport-sync -c /etc/shairport-sync.conf

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
