# recipes-multimedia/mpd/files/mpd-camel.service
[Unit]
Description=Music Player Daemon (Camel Audio - scx_bpfland scheduled)
After=network.target sound.target scx-bpfland.service

[Service]
Type=notify
User=root
Group=root

CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80

# Memory locking to prevent audio interruption from swapping
MemoryLocking=true
LimitMEMLOCK=infinity

ExecStartPre=/bin/mkdir -p /root/Music /root/.mpd/playlists /var/log
ExecStart=/usr/bin/mpd --systemd /etc/mpd-camel.conf

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=upmpdcli.service
